# py-ddns/rootfs/etc/services.d/py-ddns/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the py-ddns service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# The pyddns package reads /data/options.json itself (and installs SIGTERM/
# SIGINT handlers), so nothing to fan out from bashio::config here.
bashio::log.info "Starting py-ddns updater..."

## Run your program
exec python3 -m pyddns
